# net-misc/rpi3-ethfix/files/1.0.1/rpi3-ethfix.service
[Unit]
Description=Work around RPi3B+ Ethernet performance issues
After=network.target

[Service]
ExecStart=/usr/lib/rpi-scripts/bin/rpi3-ethfix start
Type=oneshot

[Install]
WantedBy=multi-user.target
