# x11-misc/rpi3-safecursor/files/init.d_rpi3-safecursor-1
#!/sbin/openrc-run
# Copyright (c) 2017 sakaki <sakaki@deciban.com>
# License: GPL v2 or GPL v3+
# NO WARRANTY

description="Force software cursor blitting when overscan enabled on VC4 RPi3"

depend() {
	before xdm
}

XCDIR="/usr/share/X11/xorg.conf.d"
WORKAROUND="${XCDIR}/60-swcursor.conf"

_install_or_delete_swcursor_workaround() {
	# only do this if vc4 is in use, disable_overscan=1 NOT set, and
	# the workaround file is not already present
	if [ -e "/sys/bus/platform/drivers/vc4-drm" ] && \
	grep -q "/boot" /proc/mounts && \
	! grep -q "^\s*disable_overscan=1\s*$" "/boot/config.txt" && \
	[ -d "${XCDIR}" ]; then
		if ! [ -f "${WORKAROUND}" ]; then
			einfo "You do not have disable_overscan=1 set in /boot/config.txt"
			einfo "but are using the VC4 graphics driver; as this does not"
			einfo "currently position the cursor correctly when overscan is active,"
			einfo "putting in place a workaround to enforce software blitting"
			einfo "(for more detail, see issue #17 for gentoo-on-rpi3-64bit, on GitHub)"
			cat << EOF > "${WORKAROUND}"
# This file is automatically generated by the rpi3-safecursor service
# Do not edit it!
# For more detail, see issue #17 for gentoo-on-rpi3-64bit, on GitHub.
Section "Device"
	Identifier "graphicsdriver"
	Driver "modesetting"
	Option "SWcursor" "true"
EndSection
EOF
			ewarn "Turning Xfce4 compositing off for all users"
			ewarn "This must be manually turned on again, if you want it, even"
			ewarn "if you subsequently add disable_overscan=1 to /boot/config.txt"
			find "/home" -type f -name xfwm4.xml -exec sed -i \
			's#<property name="use_compositing" type="bool" value="true"/>#<property name="use_compositing" type="bool" value="false"/>#' {} \;
			# allow user to read message; we'll be close to
			# graphical boot at this point
			sleep 5
		fi
	elif [ -f "${WORKAROUND}" ]; then
		einfo "Removing forced software cursor blitting workaround"
		rm "${WORKAROUND}"
	fi
	return 0
}

start() {
	ebegin "Starting ${SVCNAME}"
	_install_or_delete_swcursor_workaround
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	rm -f "${WORKAROUND}"
	eend 0
}
