# sys-apps/rpi-onetime-startup/files/init.d_rpi-onetime-startup-2
#!/sbin/openrc-run
# If a script /boot/startup.sh exists, run it (one time).
#
# Copyright (c) 2019 sakaki <sakaki@deciban.com>
# License: GPL v3+
# NO WARRANTY

description="Run one-time startup script, if present"
readonly SCRIPT="/boot/startup.sh"
readonly SENTINEL="/boot/autoexpand_root_partition"
readonly SENTINEL_H="/boot/autoexpand_root_partition.hold"

depend() {
	need localmount
	after local
}

_run_startup_script_if_present() {
	local RC=0
	if [[ -x "${SCRIPT}" ]]; then
		einfo "Running ${SCRIPT}..."
		${SCRIPT}
		RC=$?
	fi
	return $RC
}

_disable_this_service() {
	# don't disable during first automatic boot cycle, nor
	# during image construction
	if [[ ! -f "${SENTINEL}" && ! -f "${SENTINEL_H}" ]]; then
		einfo "Disabling service ${SVCNAME}"
		# don't continue if we can't disable (because might
		# have e.g. a read-only root)
		rc-update del ${SVCNAME} default || die "Cannot disable ${SVCNAME}"
	fi
}

start() {
   	ebegin "Starting ${SVCNAME}"
	_disable_this_service
	_run_startup_script_if_present
	eend $?
}
